# docker/gpkg-builder/make_gpkg.cron
#M  H       D   m   WD   Command
# update gpkg files 
5  2      *   *   *    /opt/etc/styles/make_gpkg.sh

